# rabbitmq/supervisord.conf
[supervisord]
nodaemon=true

[program:rabbitmq]
command=rabbitmqctl start
numprocs=1
autostart=true
autorestart=true
